# build/base/conf/libreoffice.conf
[program:libreoffice]
command=/usr/lib/libreoffice/program/soffice.bin --accept=socket,host=0,port=2002,tcpNoDelay=1;urp --headless --invisible --nocrashreport --nodefault --nologo --nofirststartwizard --norestore
autostart=true
autorestart=true
startsecs=3
user=root
numprocs=1
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0